# sv/solr/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
exec /solr/bin/solr start -force -f -p 8983 -z ${ZOOKEEPER}
